# etc/systemd/system/can-bridge.service
[Unit]
Description=Ford Fiesta MS-CAN bridge (GVRET live decoder)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/can-dashboard
Restart=on-failure
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
